# files/blork@.service
[Unit]
Description=Blork Service

[Service]
ExecStart=/path/to/blork-script.sh
StandardInput=socket
