# files/echo@.service
[Unit]
Description=Echo Service

[Service]
ExecStart=/bin/echo "Hello from echo service"
StandardInput=socket
