# files/echo.socket
[Unit]
Description=Echo Service Socket

[Socket]
ListenStream=7
Accept=yes

[Install]
WantedBy=sockets.target