# files/finger@.service
[Unit]
Description=Finger Service

[Service]
ExecStart=/usr/bin/finger
StandardInput=socket
