# files/qotd@.service
[Unit]
Description=Quote of the Day service

[Service]
ExecStart=/usr/games/fortune | /usr/games/cowsay
StandardInput=socket
